# polybot.service
[Unit]
Description=Discord Polybot Service
After=network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/MyPolyBot
ExecStart=/home/ubuntu/MyPolyBot/.venv/bin/python -m polybot.app
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/ubuntu/MyPolyBot/.env
# Uncomment and set these as needed:
#Environment=DISCORD_BOT_TOKEN=your_token_here
#Environment=YOLO_URL=http://10.0.1.90:8081/predict
#Environment=OLLAMA_URL=http://10.0.0.136:11434/api/chat
#Environment=OLLAMA_MODEL=gemma3:1b
#Environment=STATUS_SERVER_PORT=8443
# AWS Environment variables for S3 uploads
#Environment=AWS_ACCESS_KEY_ID=your_access_key_here
#Environment=AWS_SECRET_ACCESS_KEY=your_secret_key_here
#Environment=AWS_REGION=your_region_here
#Environment=AWS_DEV_S3_BUCKET=your_bucket_name_here

[Install]
WantedBy=multi-user.target 