# polybot.service
[Unit]
Description=Discord Polybot Service
After=network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/MyPolyBot
EnvironmentFile=/home/ubuntu/MyPolyBot/.env
ExecStart=/home/ubuntu/MyPolyBot/.venv/bin/python -m polybot.app
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target